# app_config/supervisord.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
